# gateway/deploy/nexus-escalation-runner.timer
[Unit]
Description=Run Nexus Escalation Runner every 2 minutes

[Timer]
OnCalendar=*-*-* *:0/2:00
AccuracySec=1s
RandomizedDelaySec=13s
Persistent=true
Unit=nexus-escalation-runner.service

[Install]
WantedBy=timers.target
