# gateway/deploy/nexus-escalation-runner.service
[Unit]
Description=Nexus Escalation Runner (per-tenant, one-shot)
After=network-online.target nexus-gateway.service
Wants=nexus-gateway.service

[Service]
Type=oneshot
User=opc
WorkingDirectory=/home/opc/afinal_gateway
EnvironmentFile=/home/opc/afinal_gateway/.env
ExecStart=/usr/bin/bash /home/opc/afinal_gateway/scripts/escalation_cron.sh
Nice=10

[Install]
WantedBy=multi-user.target
