# deploy/mac-mini/com.nexus.openclaw-operator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.openclaw-operator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/YOUR_USER/nexus-ops/openclaw/run-openclaw.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USER/nexus-ops/openclaw</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/nexus-ops/logs/openclaw.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/nexus-ops/logs/openclaw.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
